# datacapture/supervisor/conf.d/dump1090-fa.conf
[program:dump1090]
command=dump1090-fa --net --write-json /tmp/dump1090
